# data/NetworkManager/dispatcher.d/90-rncn-steamdeck-wg
#!/usr/bin/env bash
set -euo pipefail

interface="${1:-}"
action="${2:-}"
service="wg-quick@rncn-steamdeck.service"

case "$action" in
  up|vpn-up|connectivity-change|dhcp4-change|dhcp6-change)
    if systemctl is-enabled --quiet "$service"; then
      if ! systemctl is-active --quiet "$service"; then
        systemctl reset-failed "$service" || true
        systemctl start "$service" || true
      fi
    fi
    ;;
esac

exit 0
