# packaging/arch/steamos-intel-handheld.install
_steamos_intel_handheld_decky_notice() {
  local plugin_loader=/home/deck/homebrew/services/PluginLoader
  local plugin_dir=/home/deck/homebrew/plugins/steamos-intel-handheld-ec

  if [ -x "$plugin_loader" ]; then
    echo "Decky Loader detected. Charge Limit plugin files are installed at $plugin_dir."
    echo "If the panel is not visible, restart Steam or Decky Loader."
  else
    echo "Decky Loader not detected. Backend service and CLI are installed."
    echo "Steam UI Charge Limit panel requires Decky Loader; install Decky Loader first, then reinstall steamos-intel-handheld or rerun the bootstrap."
  fi

  return 0
}

_steamos_intel_handheld_gamescope_notice() {
  echo "MSI Claw 8 AI+ gamescope display profile and session hooks are installed."
  echo "Restart the gamescope session or reboot for the display profile and native-panel wrapper to take effect."
  return 0
}

_steamos_intel_handheld_enable_services() {
  if ! command -v systemctl >/dev/null 2>&1; then
    return 0
  fi

  systemctl daemon-reload || true
  systemctl enable steamos-intel-handheld-restore.service || true
  systemctl enable steamos-intel-handheld-power-control.service || true
  systemctl disable steamos-intel-handheld-steamos-manager-remote.service 2>/dev/null || true
  systemctl stop steamos-intel-handheld-power-control.service || true
  busctl call org.freedesktop.DBus /org/freedesktop/DBus org.freedesktop.DBus ReloadConfig || true

  if [ -x /usr/bin/steamos-intel-handheld-restore-etc ]; then
    /usr/bin/steamos-intel-handheld-restore-etc --apply || true
  fi

  _steamos_intel_handheld_restart_user_manager_without_provider || true
  systemctl start steamos-intel-handheld-power-control.service || true

  return 0
}

_steamos_intel_handheld_restart_user_manager_without_provider() {
  local uid runtime_dir

  uid="$(id -u deck 2>/dev/null || true)"
  if [ -z "$uid" ]; then
    return 0
  fi

  runtime_dir="/run/user/$uid"
  if [ ! -S "$runtime_dir/bus" ]; then
    return 0
  fi

  timeout 45 runuser -u deck -- env \
    XDG_RUNTIME_DIR="$runtime_dir" \
    DBUS_SESSION_BUS_ADDRESS="unix:path=$runtime_dir/bus" \
    systemctl --user restart steamos-manager.service || true
}

post_install() {
  _steamos_intel_handheld_enable_services || true
  _steamos_intel_handheld_decky_notice || true
  _steamos_intel_handheld_gamescope_notice || true
}

post_upgrade() {
  _steamos_intel_handheld_enable_services || true
  _steamos_intel_handheld_decky_notice || true
  _steamos_intel_handheld_gamescope_notice || true
}
