# data/systemd/user/gamescope-mangoapp.service.d/10-rivoreo-mangoapp.conf
[Service]
ExecStart=
ExecStart=/opt/steamos-intel-handheld/bin/mangoapp
