# data/systemd/user/gamescope-session.service.d/20-native-panel-resolution.conf
[Service]
Environment=PATH=/opt/steamos-intel-handheld/bin:/usr/local/sbin:/usr/local/bin:/usr/bin:/bin
